# auditor/crontab
0 * * * * bash /scripts/audit.sh >> /reports/cron.log 2>&1
